# youtube-offline.service
[Unit]
Description=YouTube Offline - Automatic Playlist Downloader
After=network.target

[Service]
Type=simple
WorkingDirectory=WORKING_DIR_PLACEHOLDER
ExecStart=NODE_PATH_PLACEHOLDER WORKING_DIR_PLACEHOLDER/dist/server.js
Restart=always
RestartSec=10

# Logging - use journald (view with: journalctl --user -u youtube-offline -f)
StandardOutput=journal
StandardError=journal

# Environment
Environment="NODE_ENV=production"
Environment="PATH=PATH_PLACEHOLDER:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
